# root/etc/s6-overlay/s6-rc.d/init-nextcloud-aliases/up
/etc/s6-overlay/s6-rc.d/init-nextcloud-aliases/run
